# scripts/com.missioncontrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.missioncontrol</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mac/Downloads/Claw/mission-control/scripts/launch-mission-control.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/mac/Downloads/Claw/mission-control</string>
    <key>StandardOutPath</key>
    <string>/tmp/mission-control.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mission-control.stderr.log</string>
</dict>
</plist>
